# build/supervisor_all.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf  -g 'daemon off;'
;process_name=%(program_name)s_%(process_num)02d

;defaul unexpected, restart only if exited abnormally
autorestart=false
;do not wait before transition to running
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
;do not rotate log file
stdout_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm7.2 -R -F -c /etc/php/7.2/fpm/php-fpm.conf
;process_name=%(program_name)s_%(process_num)02d
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:ssh]
command=/usr/sbin/sshd -D
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:supercron]
command=php /local/update_ssh_keys.php
autorestart=true
startsecs=0
startretries=0
user=git
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
